# .config/systemd/user/redshift.service.d/display.conf
[Service]
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 4
